# supervisord-ssh.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
# deploy 用户通过组权限访问 supervisorctl
chown=root:deploy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd-error.log

[program:llm-http-proxy]
command=/opt/proxy/llm-http-proxy -addr :8080 -persist /data/stats.json
autostart=true
autorestart=true
stdout_logfile=/public/logs/proxy.log
stderr_logfile=/public/logs/proxy-error.log
